# k8s/config/yapi.cm.supervisord/supervisord.server.ini
[program:server]
command=sh /data/server_start.sh
autostart=true
autorestart=false
stderr_logfile=/var/log/supervisor/supervisord.log
stdout_logfile=/var/log/supervisor/supervisord.log
startsecs=0
